# linux/crontab
*/5 * * * * curl https://api.repo.com.ar/cron/categorias_productos
